# systemd/avr300osc.service
[Unit]
Description=AVR300 OSC Controller

[Service]
Type=simple
ExecStart=/usr/local/bin/avr300osc

[Install]
WantedBy=multi-user.target
